# supervisord.conf
# ============================================
# Sora v2.0 — Supervisord Configuration
#
# 프로세스 구조:
#   1. redis       — 메시지 큐 + pubsub
#   2. gateway     — FastAPI (HTTP/WS, PC Agent Hub, TG webhook)
#   3. brain       — LLM Worker (Redis 큐 소비, 도구 실행)
#   4. scheduler   — APScheduler (크론잡)
#   5. cloudflared — HTTPS 터널 (선택)
#
# 핵심 원칙: gateway는 LLM을 직접 호출하지 않음
# ============================================

[supervisord]
nodaemon=true
user=root
logfile=/app/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info

# ── 1. Redis (메시지 버스) ──
[program:redis]
command=redis-server --save "" --appendonly no --maxmemory 128mb --maxmemory-policy allkeys-lru
autostart=true
autorestart=true
startretries=3
startsecs=2
priority=1
stdout_logfile=/app/logs/redis.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/app/logs/redis_err.log
stderr_logfile_maxbytes=5MB

# ── 2. Gateway (FastAPI — 절대 블로킹 없음) ──
[program:gateway]
command=python -m uvicorn src.core.sora_dashboard:app --host 0.0.0.0 --port 7700 --log-level warning
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=3
priority=10
environment=SORA_CLOUD_MODE="true"
stdout_logfile=/app/logs/sora_dashboard.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/app/logs/sora_dashboard_err.log
stderr_logfile_maxbytes=10MB

# ── 3. Brain Worker (LLM — Gateway 안정화 후 시작) ──
[program:brain]
command=python -m src.core.brain.worker
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=5
priority=50
stdout_logfile=/app/logs/brain.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/app/logs/brain_err.log
stderr_logfile_maxbytes=10MB

# ── 4. Scheduler (크론잡) ──
[program:scheduler]
command=python neo_genesis_daemon.py
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=10
priority=30
stdout_logfile=/app/logs/daemon_stdout.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/app/logs/daemon_stderr.log
stderr_logfile_maxbytes=10MB

# ── 5. Cloudflare Tunnel (선택) ──
[program:cloudflared]
command=/bin/sh -c 'if [ -f /app/secrets/cloudflared-token.txt ]; then cloudflared tunnel run --token "$(cat /app/secrets/cloudflared-token.txt)"; else cloudflared tunnel --config /app/secrets/cloudflared-config.yml run; fi'
autostart=%(ENV_CF_TUNNEL_ENABLED)s
autorestart=true
startretries=10
startsecs=5
priority=5
stdout_logfile=/app/logs/cloudflare_tunnel.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/app/logs/cloudflare_tunnel_err.log
stderr_logfile_maxbytes=5MB
